# etc/systemd/system/centreon-centacl.timer
[Unit]
Description=Run Centreon CentACL every minute

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=1s

[Install]
WantedBy=timers.target
