# etc/systemd/system/centreon-log-analyser.timer
[Unit]
Description=Run Centreon log analyser every minute

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=1s

[Install]
WantedBy=timers.target
